# fan/ilo-fan.timer
[Unit]
Description=Set iLO fan after reboot

[Timer]
OnBootSec=30s

[Install]
WantedBy=timers.target
